# ansible/roles/uploader/files/etc/systemd/system/uploader.service
[Unit]
Description=Piazza file upload form
DefaultDependencies=no
After=
Conflicts=
Before=
Requires=

[Service]
ExecStart=/opt/uploader/bin/uploader
## To not use augenrules, copy this file to /etc/systemd/system/auditd.service
## and comment/delete the next line and uncomment the auditctl line.
## NOTE: augenrules expect any rules to be added to /etc/audit/rules.d/
# ExecStartPost=
# ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multiuser.target
